# tc_radarr/sv/radarr/finish
#!/bin/bash
## headphones/finish

echo "[APPLICATION] Radarr shutting down."

kill -9 $(ps ax | grep NzbDrone.exe | awk '{print $1}')
wait $!
